# templates/gunicorn.conf.j2
[program:gunicorn]
command=/root/venvs/env/scripts/gunicorn.sh

