# templates/supervisor.conf.j2
[program:virtualenv]
command=/root/venvs/env/scripts/runcms.sh


